# pagerduty/meta.cron.hourly
#!/bin/bash

PROBE=`readlink -f "$0"`
PROBE=`dirname "$PROBE"`/probe

$PROBE z.arrowfs.org ping z.arrowfs.org
