# puppetmaster/modules/bootstrap/files/bootstrap-startup.conf
# Bootstrap machine-startup tasks

# This is a task (not a service) that ensures we run the bootstrap
# startup tasks on boot, after making sure nginx and config have
# reported they started.

description "aerofs bootstrap tasks run on startup"
author "AeroFS <team@aerofs.com>"

start on startup and started nginx and started config

task

chdir /opt/bootstrap
exec /usr/bin/aerofs-bootstrap-taskfile /opt/bootstrap/tasks/startup.tasks
