# packaging/synctime/etc/init/synctime.conf
# synctime server script
# this is completely copied from devman

description "synctime server"
author "Mark Jeffrey <mark@aerofs.com"

# When to start the service
start on runlevel [2345]

# When to stop the service
stop on runlevel [016]

# Automatically restart process if crashed
respawn

# Start the process and outputs only system err to synctime.err.log
# TODO(MJ) devman formerly had --chuid devman, but there currently isn't a
# synctime user on the tts box
exec start-stop-daemon --start --chdir /opt/synctime --exec /opt/synctime/synctime 2>> /var/log/aerofs/synctime.err.log
