# packaging/vmhost/etc/init/aerofs-vmhost-bootstrap.conf
# aerofs-vmhost-boostrap - AeroFS bootstrap web application.

description "AeroFS bootstrap web application."

start on runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/local/bin/pserve /usr/local/aerofs/python/bootstrap/config.ini
