# puppetmaster/modules/rocklog/files/gunicorn.conf
description "gunicorn"

start on runlevel [2]
stop on runlevel [016]

chdir /opt/rocklog

exec gunicorn rocklog:app
