# puppetmaster/modules/rocklog/files/tornado.conf
description "tornado rocklog backend"

start on [2345]
stop on [016]

respawn

chdir /opt/rocklog
exec python /opt/rocklog/tornado_server.py
